# cookbooks/consul_server/files/default/supervisord_vault.conf
[program:vault]
command=/usr/local/bin/vault server -config /etc/vault/vault.conf
autostart=true
autorestart=true
stderr_logfile=/var/log/vault.err.log
stdout_logfile=/var/log/vault.out.lo